# files/ddclient/etc/cron.daily/ddclient-update
#!/bin/sh -e
# Update container to latest
podman pull lscr.io/linuxserver/ddclient:latest

# Restart pod to use latest container
systemctl restart pod-ddclient.service
